# tests/mezzver/debian-mezz588-gs399/oblong-mezzanine-init5.88.postrm
#!/bin/sh
# postrm script for mezzanine common scripts
#
# see: dh_installdeb(1)

set -e
set -x

# summary of how this script can be called:
#         * <postrm> `remove'
#         * <postrm> `purge'
#         * <old-postrm> `upgrade' <new-version>
#         * <new-postrm> `failed-upgrade' <old-version>
#         * <new-postrm> `abort-install'
#         * <new-postrm> `abort-install' <old-version>
#         * <new-postrm> `abort-upgrade' <old-version>
#         * <disappearer's-postrm> `disappear' <r>overwrit>r> <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

MAJOR=$( lsb_release -sr | cut -d"." -f1 )
MAJOR_WITH_SYSTEMD=15

remove_file ()
{
  if [ -e $1 ] || [ -h $1 ]
  then
    rm $1
  fi
}

case "$1" in
  purge|remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        remove_file /etc/logrotate.d/ob-rotate-mezzanine

        # because x11 may have overwritten this file
        if [ -e /etc/X11/default-display-manager.obbak ]
        then
          mv -f /etc/X11/default-display-manager.obbak /etc/X11/default-display-manager
        fi

        if [ ${MAJOR} -lt ${MAJOR_WITH_SYSTEMD} ]
        then
          # 12.04 only
          remove_file /etc/init.d/mz-pools-setup
          remove_file /etc/init.d/ob-http-ctl
          remove_file /etc/init.d/check-mezz
          remove_file /etc/init.d/check-firmware
          remove_file /etc/init.d/mz-plasma-web-proxy
          remove_file /etc/init.d/x11
        else
          # 16.04 only
          rm -f /etc/systemd/journald.conf.d/oblong-persistence.conf || true
        fi

        remove_file /etc/oblong/startup-commons
        remove_file /etc/oblong/xinitrc
        remove_file /etc/nginx/sites-enabled/obsslproxy
        remove_file /etc/network/if-pre-up.d/mz-firewall-up

        # remove the line we added to the nginx conf file in the postinst script
        NGINX_FILE="/etc/default/nginx"
        if [ -f ${NGINX_FILE} ]; then
          sed -i '\/etc\/oblong\/nginx.conf/d' ${NGINX_FILE}
        fi

        if [ -f /etc/monit/monitrc.orig.mz ]
        then
          mv /etc/monit/monitrc.orig.mz /etc/monit/monitrc
        fi

        # and remove the 16.04 version of the above nginx fix
        if [ ${MAJOR} -ge ${MAJOR_WITH_SYSTEMD} ]
        then
          # the 16.04 systemd nginx.service does not pay attention to
          # the above modified file, so...
          # remove oblong from nginx systemd unit (bug 16068)
          lx=$( readlink --canonicalize /lib/systemd/system/nginx.service )
          if [ -f "$lx" ]
          then
            if [ $( grep -c oblong $lx ) -gt 0 ]
            then
              sed -i -e "s/-c \/etc\/oblong\/nginx.conf //" $lx || true
              systemctl daemon-reload || true
            fi
          fi
        fi

        # look through /etc/sysctl.d/*.conf
        if [ ${MAJOR} -lt ${MAJOR_WITH_SYSTEMD} ]
        then
          # 12.04 only
          invoke-rc.d procps start
        fi

        # edit fstab to remove the ramdisk for pools we created
        FS_TAB="/etc/fstab"
        if [ -f ${FS_TAB} ]; then
          sed -i '/automatically created by oblong-mezzanine-init.postinst/d' ${FS_TAB}
          sed -i '/\/mnt\/poolsramdisk/d' ${FS_TAB}
        fi

        # Remove check-firmware
        if [ ${MAJOR} -lt ${MAJOR_WITH_SYSTEMD} ]
        then
          # 12.04 only
          update-rc.d -f check-firmware remove
        fi
  ;;

  *)
    echo "postrm called with unknown argument \`$1'" >&2
    exit 1
  ;;

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
