# tests/mezzver/debian-mezz588-gs399/oblong-mezzanine5.88.postinst
#!/bin/sh
# postinst script for mezzanine
#
# see: dh_installdeb(1)

set -e
set -x

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

MAJOR=$( lsb_release -sr | cut -d"." -f1 )
MAJOR_WITH_SYSTEMD=15

create_link ()
{
  rm -f $2 # all init scripts land in /etc/init.d/mz now (16.04 port)
  ln -s $1 $2
}

case "$1" in
    configure)
        if [ ${MAJOR} -lt ${MAJOR_WITH_SYSTEMD} ]
        then
          # links are only made in 12.04, systemd in 16.04 is
          # confused by the mere existence of sysv scripts
          for f in mezzanine \
                   siemcy \
                   pygiandros \
                   thumbhammer \
                   rtsp-remote-provisioner-spawner \
                   giecy \
                   qm-provisioner-bm \
                   qm-merchant-inogeni \
                   rtsp-viddle-server \
                   rtsp-inroom-server \
                   ob-cthulhu \
                   remote-participation-pool-server \
                   meetman \
                   iperf \
                   vtci
          do
            create_link /etc/init.d/mz/$f /etc/init.d/$f
          done
          update-rc.d -f ondemand disable || true
        else
          # in 16.04 monit is superseded by systemd
          systemctl disable monit || true
          systemctl disable ondemand || true
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
